# services/startcar.service
[Unit]
Description=Start the Console Display
[Service]
ExecStart=/home/brian/start.sh 